# conf/gunicorn.service
[Unit]
Description=gunicorn daemon
After=network.target

[Service]
User=www-data
Group=www-data
WorkingDirectory=/app/CocktailManager
ExecStart=/usr/local/bin/gunicorn --workers 3 --bind unix:/app/CocktailManager/cocktail.sock CocktailManager.wsgi:application

[Install]
WantedBy=multi-user.target