# conf/rc.local
#!/bin/sh -e
/usr/local/bin/uwsgi --ini /data/CocktailManager/conf/CM-wsgi.ini
service nginx start
service postgersql start
python3 /data/CocktailManager/CocktailManager/manage.py makemigrations
python3 /data/CocktailManager/CocktailManager/manage.py migrate
echo yes | python3 /data/CocktailManager/CocktailManager/manage.py collectstatic
exit 0
